# 1.20-ubuntu-bionic-s6/etc/services.d/confd/run
#!/usr/bin/execlineb -P

with-contenv
confd -node http://etcd:4001 -config-file /etc/confd/confd.toml
