# 1.5/config/etc/supervisor/conf.d/confd.conf
[program:confd]
command=sh -c "exec confd -node http://$ETCD_PORT_4001_TCP_ADDR:$ETCD_PORT_4001_TCP_PORT -config-file /etc/confd/confd.toml"
redirect_stderr=true
stopsignal=TERM
autorestart=true
