# config/etc/supervisor/conf.d/confd.conf
[program:confd]
command=confd -node http://etcd:4001 -config-file /etc/confd/confd.toml
redirect_stderr=true
stopsignal=TERM
autorestart=true
